# main/data/dinox.service
[Unit]
Description=DinoX — XMPP Instant Messenger
Documentation=https://github.com/rallep71/dinox
# This unit is for manual start/stop/restart via systemctl --user.
# Autostart on login is handled by XDG autostart:
#   ~/.config/autostart/im.github.rallep71.DinoX.desktop
# Toggle it with the "Start Automatically" switch in DinoX preferences.
PartOf=graphical-session.target
After=graphical-session.target
StartLimitBurst=3
StartLimitIntervalSec=30s

[Service]
Type=simple
ExecStart=/usr/bin/dinox
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=graphical-session.target
